# off-button.service
[Unit]
Description=Off button service

[Service]
Type=simple
ExecStart=/usr/sbin/off-button

[Install]
WantedBy=multi-user.target


